# f5-tts-api.service
[Unit]
Description=F5-TTS OpenAI-Compatible API (Russian model)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=drynw
WorkingDirectory=/home/drynw/models/f5-tts
Environment=PATH=/home/drynw/f5-tts-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/drynw/f5-tts-env/bin/python3 /home/drynw/models/f5-tts/api_server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
