# f5-tts-web.service
[Unit]
Description=F5-TTS Character Manager Web UI
After=network-online.target f5-tts-api.service
Wants=network-online.target f5-tts-api.service

[Service]
Type=simple
User=drynw
WorkingDirectory=/home/drynw/models/f5-tts
Environment=PATH=/home/drynw/f5-tts-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=F5_TTS_BASE_DIR=/home/drynw/models/f5-tts
Environment=F5_TTS_VENV=/home/drynw/f5-tts-env
Environment=F5_TTS_API_URL=http://localhost:8000
ExecStart=/home/drynw/f5-tts-env/bin/python3 /home/drynw/models/f5-tts/web_ui.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
